# scripts/etc/systemd/system/qmanager.target
# /etc/systemd/system/qmanager.target
[Unit]
Description=QManager Services
After=multi-user.target

[Install]
WantedBy=multi-user.target
