# scripts/etc/systemd/system/qmanager-eth-link.service
# /etc/systemd/system/qmanager-eth-link.service
[Unit]
Description=QManager Ethernet Link Speed Persistence
After=network.target qmanager-setup.service
PartOf=qmanager.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sh -c 'command -v ethtool >/dev/null 2>&1 && [ -d /sys/class/net/eth0 ]'
ExecStart=/usr/bin/qmanager_eth_link_apply
ExecStop=/usr/bin/qmanager_eth_link_reset

[Install]
WantedBy=qmanager.target
