# scripts/etc/systemd/system/qmanager-bandwidth-monitor.service
# /etc/systemd/system/qmanager-bandwidth-monitor.service
[Unit]
Description=QManager Bridge Traffic Monitor Binary
After=qmanager-bandwidth-genconf.service qmanager-bandwidth-websocat.service
Requires=qmanager-bandwidth-genconf.service
PartOf=qmanager-bandwidth.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /tmp/quecmanager
ExecStart=/usr/bin/bridge_traffic_monitor_rm551
ExecStopPost=/bin/sh -c 'rm -f /tmp/quecmanager/bridge_traffic_monitor /tmp/quecmanager/bridge_traffic_monitor.pid'
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=qmanager-bandwidth.target
